# packaging/deb/heimdalld/DEBIAN/postinst
#!/bin/bash
# This is a postinstallation script so the service can be configured and started when requested
#
sudo systemctl enable heimdalld.service
sudo systemctl daemon-reload
